# etc/s6-overlay/s6-rc.d/copy-files/up
#!/command/execlineb
foreground { bash -e /usr/local/bin/pihole-unbound/copy-files.sh }