# etc/s6-overlay/s6-rc.d/setup-unbound-log/up
#!/command/execlineb
foreground { bash -e /usr/local/bin/pihole-unbound/setup-unbound-log.sh }