# etc/s6-overlay/s6-rc.d/update-root-hints/up
#!/command/execlineb
foreground { bash -e /usr/local/bin/pihole-unbound/update-root-hints.sh }