# etc/s6-overlay/s6-rc.d/update-root-hints/up
#!/usr/bin/execlineb -P
with-contenv bash -c '
ROOT_HINTS_FILE="/var/lib/unbound/root.hints"

echo "Updating root.hints file..."
if curl -o "$ROOT_HINTS_FILE" https://www.internic.net/domain/named.cache; then
    echo "Successfully updated root.hints file."
else
    echo "Failed to update root.hints file."
    exit 1
fi
'