# etc/s6-overlay/s6-rc.d/add-local-dns/up
#!/usr/bin/execlineb -P
with-contenv bash -c '
LOCALDNS_FILE="/etc/pihole/custom.list"
: > "$LOCALDNS_FILE"

IFS=","
for entry in $LOCALDNS; do
    ip="${entry%%=*}"
    domain="${entry#*=}"

    if ! grep -Fxq "$ip $domain" "$LOCALDNS_FILE"; then
        echo "$ip $domain" >> "$LOCALDNS_FILE"
    fi
done

echo "Local DNS entries update complete."
'