# etc/s6-overlay/s6-rc.d/unbound-log-setup/up
#!/usr/bin/execlineb -P
with-contenv bash -c '
mkdir -p /var/log/unbound
touch /var/log/unbound/unbound.log
chown unbound /var/log/unbound/unbound.log
echo "Unbound log setup complete."
'