# etc/s6-overlay/s6-rc.d/restart-dns/up
#!/usr/bin/execlineb -P
with-contenv bash -c '
set -e
echo "Restarting Pi-hole DNS service..."
pihole restartdns
'