# etc/s6-overlay/s6-rc.d/add-cname-records/up
#!/usr/bin/execlineb -P
with-contenv bash -c '
CNAME_FILE="/etc/dnsmasq.d/05-cname.conf"
: > "$CNAME_FILE"

IFS=","
for entry in $CNAME_ENTRIES; do
    alias="${entry%%=*}"
    target="${entry#*=}"

    if ! grep -Fxq "cname=$alias,$target" "$CNAME_FILE"; then
        echo "cname=$alias,$target" >> "$CNAME_FILE"
    fi
done

echo "CNAME records update complete."
'