# etc/s6-overlay/s6-rc.d/add-cname-records/up
#!/command/execlineb
foreground { bash -e /usr/local/bin/add-cname-records.sh }